# systemd/getty@tty1.service.d/skip-username.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -- rasul' --noclear --skip-login - $TERM
